# systemd/claude-discord-bridge.service
[Unit]
Description=Claude-Discord Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/victor/agent-env/workspace/claude-discord-bridge
ExecStart=/home/victor/agent-env/workspace/claude-discord-bridge/bin/vai
ExecStop=/home/victor/agent-env/workspace/claude-discord-bridge/bin/vexit
Environment=HOME=/home/victor
Environment=PATH=/home/victor/agent-env/workspace/claude-discord-bridge/bin:/home/victor/.local/bin:/usr/local/bin:/usr/bin:/bin
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=default.target
